# launchd/com.chrisren.autonomy-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.chrisren.autonomy-sweep</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="$HOME/.claude/bin:$PATH"; ~/.claude/scripts/autonomy-sweep.sh</string>
  </array>
  <key>StartInterval</key>    <integer>300</integer>
  <key>RunAtLoad</key>        <false/>
  <key>ProcessType</key>      <string>Background</string>
  <key>Nice</key>             <integer>5</integer>
  <key>StandardOutPath</key>  <string>/Users/chrisren/.claude/logs/autonomy-sweep.out.log</string>
  <key>StandardErrorPath</key><string>/Users/chrisren/.claude/logs/autonomy-sweep.err.log</string>
</dict>
</plist>
