# launchd/com.chrisren.watch-getAppState-fix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chrisren.watch-getAppState-fix</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>source $HOME/.zshrc 2>/dev/null; PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin $HOME/.claude/scripts/watch-getAppState-fix.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>7</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/chrisren/.claude/logs/launchd-watch-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chrisren/.claude/logs/launchd-watch-stderr.log</string>
</dict>
</plist>
